# config/xboxdrv/xboxdrv.service
[Unit]
Description=xboxdrv controller remaps
After=rc-local.service

[Install]
WantedBy=multi-user.target
