# config/powerkit/powerkit.service
[Unit]
Description=Hardware case safe shutdown
After=rc-local.service

[Service]
Type=simple
ExecStart=python3 /opt/retropie/supplementary/powerkit/cli.py /opt/retropie/configs/all/powerkit.cfg

[Install]
WantedBy=multi-user.target
